# src/main/etc/install/tomcat.service
[Unit]
Description=Tomcat webserver
After=syslog.target network.target

[Service]
Type=forking
User=root
ExecStart=/usr/local/share/tomcat/bin/startup.sh 
ExecStop=/usr/local/share/shutdown.sh

[Install]
WantedBy=multi-user.target
